# Composer/dailyCI/task.cron
30 3 * * * /bin/sh /home/azureuser/LuisE2E/runAndSendMail.sh >> ./e2etest.log 2>&1
